# net-misc/downloaddaemon/files/downloadd_openrc
#!/sbin/runscript
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need net
}

start() {

	ebegin "Starting DownloadDaemon"
	start-stop-daemon --start --quiet --exec "/usr/bin/DownloadDaemon" --user downloadd -- -d -p /run/downloadd.pid
	eend $?

}

stop() {

	ebegin "Stopping DownloadDaemon"
	start-stop-daemon --stop --quiet --pidfile "/run/downloadd.pid"
	eend $?

}
